# tradetherapybot.service
[Unit]
Description=Trade Therapy Bot - Telegram bot for subscription management
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/TradeTherapyBot
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /path/to/TradeTherapyBot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

