# xxx/kopiya-iksuyemsya.service
[Unit]
Description=Копия иксуюемся - Telegram bot (userbot + aiogram)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cursor_test/xxx
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /opt/cursor_test/xxx/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
